# systemd/node_exporter.service
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
Wants=network-online.target
After=network-online.target

[Service]
User=node_exporter
Group=node_exporter
Type=simple
Restart=on-failure
RestartSec=5s
TimeoutStartSec=0
TimeoutStopSec=30s

ExecStart=/opt/node_exporter/node_exporter \
  --web.listen-address=":9100" \
  --path.procfs=/proc \
  --path.sysfs=/sys \
  --path.rootfs=/ \
  --collector.filesystem.mount-points-exclude="^/(sys|proc|dev|host|etc|rootfs/var/lib/docker/containers|rootfs/var/lib/docker/overlay2|rootfs/run/docker/netns|rootfs/var/lib/docker/aufs)($$|/)"
\
  --collector.filesystem.fs-types-exclude="^(autofs|binfmt_misc|bpf|cgroup2?|configfs|debugfs|devpts|devtmpfs|fusectl|hugetlbfs|iso9660|mqueue|nsfs|overlay|proc|procfs|pstore|rpc_pipefs|securityfs|selinuxfs|squashfs|sysfs|tracefs)$$" \
  --collector.netdev.device-exclude="^(veth.*|docker.*|br-.*|lo)$$" \
  --collector.diskstats.ignored-devices="^(ram|loop|fd|nbd|sr|dm-|md)[0-9]+$$" \
  --no-collector.arp \
  --no-collector.bcache \
  --no-collector.bonding \
  --no-collector.conntrack \
  --no-collector.edac \
  --no-collector.entropy \
  --no-collector.hwmon \
  --no-collector.infiniband \
  --no-collector.ipvs \
  --no-collector.mdadm \
  --no-collector.nfs \
  --no-collector.nfsd \
  --no-collector.powersupplyclass \
  --no-collector.rapl \
  --no-collector.thermal_zone \
  --no-collector.time \
  --no-collector.timex \
  --no-collector.wifi \
  --no-collector.xfs \
  --no-collector.zfs

ExecReload=/bin/kill -HUP $$MAINPID
KillMode=mixed

NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateTmp=true

LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
